# systemd/stremio-lan-remote-watchdog.service
[Unit]
Description=Stremio LAN Remote shell watchdog
After=network-online.target

[Service]
Type=oneshot
ExecStart=%h/dev/stremio-lan-remote/scripts/watchdog.sh
